# etc/cron.d/powerdnsadmin
# min   hour    day     month   weekday   user      command
*/5     *       *       *       *         gunicorn  /usr/local/sbin/powerdnsadmin-update
